# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── 1. MariaDB ────────────────────────────────────────────────────────────────
[program:mysql]
command=mysqld --user=mysql --datadir=/var/lib/mysql
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 2. Backend NestJS ─────────────────────────────────────────────────────────
; Espera a MySQL por TCP con el usuario folklosoft, luego migra y arranca
[program:backend]
command=sh -c "until mariadb -h 127.0.0.1 -P 3306 -u folklosoft -p\"$DB_PASSWORD\" folklosoft -e 'SELECT 1' 2>/dev/null; do echo 'backend: esperando MySQL...'; sleep 2; done && npx prisma db push --accept-data-loss && node dist/prisma/seed.js && node dist/src/main.js"
directory=/app/backend
autostart=true
autorestart=true
priority=2
environment=NODE_ENV="production",PORT="9001",DATABASE_URL="%(ENV_DATABASE_URL)s",JWT_SECRET="%(ENV_JWT_SECRET)s",FRONTEND_URL="%(ENV_FRONTEND_URL)s",JWT_ACCESS_EXPIRES_IN="%(ENV_JWT_ACCESS_EXPIRES_IN)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 3. Adminer (gestor web de BD) ─────────────────────────────────────────────
[program:adminer]
command=php83 -S 0.0.0.0:9004 /app/adminer/index.php
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 4. Frontend Next.js ───────────────────────────────────────────────────────
[program:frontend]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=3
environment=NODE_ENV="production",PORT="9002",BACKEND_URL="http://127.0.0.1:9001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
