# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── 1. PostgreSQL (solo cuando NO hay DB externa) ─────────────────────────────
[program:postgres]
command=/usr/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=%(ENV_START_DB)s
autorestart=%(ENV_START_DB)s
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 2. Backend NestJS ─────────────────────────────────────────────────────────
[program:backend]
command=/app/backend/start-backend.sh
directory=/app/backend
autostart=true
autorestart=true
priority=2
environment=NODE_ENV="production",PORT="4002",DATABASE_URL="%(ENV_DATABASE_URL)s",JWT_SECRET="%(ENV_JWT_SECRET)s",FRONTEND_URL="%(ENV_FRONTEND_URL)s",JWT_ACCESS_EXPIRES_IN="%(ENV_JWT_ACCESS_EXPIRES_IN)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 3. Adminer (gestor web de BD) ─────────────────────────────────────────────
[program:adminer]
command=php83 -S 0.0.0.0:9004 /app/adminer/index.php
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 4. Frontend Next.js ───────────────────────────────────────────────────────
[program:frontend]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
priority=3
environment=NODE_ENV="production",PORT="3002",BACKEND_URL="http://127.0.0.1:4002"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
